# infrastructure/systemd/aruna-node.service
[Unit]
Description=ARUNA Network Core Node
After=network.target

[Service]
Type=simple
User=aruna
Group=aruna
WorkingDirectory=/var/lib/aruna-node
ExecStart=/usr/local/bin/aruna-node daemon --p2p-port 9000 --rpc-port 8080
Restart=always
RestartSec=5s
LimitNOFILE=65535

# Security Hardening Sandbox Properties
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/aruna-node
PrivateTmp=true
ProtectKernelTunables=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target
